# deployment/v0.9.2/ollama_output/fixes/fix13_new_d3kos-settings-api.service
[Unit]
Description=d3kOS Settings Action API
After=network.target

[Service]
Type=simple
User=d3kos
WorkingDirectory=/opt/d3kos/services
ExecStart=/usr/bin/python3 /opt/d3kos/services/settings-api.py
Restart=on-failure

[Install]
WantedBy=multi-user.target